# options-dashboard-api.service
[Unit]
Description=Options Dashboard FastAPI Backend
After=network.target
Wants=options-dashboard-nginx.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/options_dashboard
Environment="PATH=/home/ubuntu/options_dashboard/.venv/bin"
ExecStart=/home/ubuntu/options_dashboard/.venv/bin/python -m uvicorn api.main:app --host 127.0.0.1 --port 8000 --workers 2
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=process
Restart=always
RestartSec=10

# Security
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/options_dashboard/data /home/ubuntu/options_dashboard/options_data /home/ubuntu/options_dashboard/logs

[Install]
WantedBy=multi-user.target
